# sigma_elastalert_rules/macos/file_event/file_event_macos_startup_items.yaml
alert:
- debug
description: Detects creation of startup item plist files that automatically get executed
  at boot initialization to establish persistence.
filter:
- query:
    query_string:
      query: (file_name.keyword:*\/Library\/StartupItems\/* OR file_name.keyword:*.plist)
index: wazuh-alerts*
name: file_event_macos_startup_items
priority: 4
realert:
  minutes: 0
type: any

is_enabled: false
timestamp_type: '@timestamp'
